# lifeshield.conf
# Run script(s) to monitor LifeShield security system

description     "Monitor LifeShield"

start on (filesystem)
stop on starting rc RUNLEVEL=[016]
respawn

exec php /PATH/TO/lifeshield.php

